# deploy/letter_counter.ini
[program:letter_counter]
user=marc
command=/home/marc/bin/python /home/marc/python_test_apps/letter_counter/server.py
directory=/home/marc/python_test_apps/letter_counter/
autostart=true
autorestart=true
logfile=log.txt
